# packaging/systemd/rcod.user.service
[Unit]
Description=RCOD - Remote Control On Demand (user service)

[Service]
Type=simple
WorkingDirectory=<STATE_DIR>
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:%h/.local/bin
ExecStart=<BIN_PATH> --config <CONFIG_PATH> --state-dir <STATE_DIR>
Restart=on-failure
RestartSec=5s
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
