# packaging/launchd/ai.zevro.rcod.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.zevro.rcod.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>@BIN_PATH@</string>
    <string>--config</string>
    <string>@CONFIG_PATH@</string>
    <string>--state-dir</string>
    <string>@STATE_DIR@</string>
  </array>
  <key>WorkingDirectory</key>
  <string>@STATE_DIR@</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:@USER_HOME@/.local/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ThrottleInterval</key>
  <integer>5</integer>
  <key>StandardOutPath</key>
  <string>@LOG_DIR@/rcod.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>@LOG_DIR@/rcod.stderr.log</string>
</dict>
</plist>
