# packaging/launchd/ai.zevro.rcod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.zevro.rcod</string>
  <key>ProgramArguments</key>
  <array>
    <string><REPO_PATH>/rcod</string>
    <string>--config</string>
    <string><CONFIG_PATH></string>
  </array>
  <key>WorkingDirectory</key>
  <string><REPO_PATH></string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string><LOG_PATH>/rcod.stdout.log</string>
  <key>StandardErrorPath</key>
  <string><LOG_PATH>/rcod.stderr.log</string>
</dict>
</plist>
